# daemon_node/tools/ros_daemon_node.service
[Unit]
Description=ros_daemon_service.

[Service]
ExecStart=/bin/bash /home/hugoliu/github/catkin_ws/src/VIW-Fusion/daemon_node/tools/run_daemon_node.sh

[Install]
WantedBy=multi-user.target